# scripts/vm-bootstrap/systemd/docker-compose-hermes.service
[Unit]
Description=Hermes docker-compose stack
After=docker.service hermes-secrets.service
Requires=docker.service hermes-secrets.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/hermes/bootstrap
EnvironmentFile=/run/hermes/env/*
ExecStartPre=/usr/bin/docker compose -f docker-compose.yml -f docker-compose.gcp.override.yml pull
ExecStart=/usr/bin/docker compose -f docker-compose.yml -f docker-compose.gcp.override.yml up -d
ExecStop=/usr/bin/docker compose -f docker-compose.yml -f docker-compose.gcp.override.yml down
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
