# scripts/vm-bootstrap/systemd/hermes-secrets.service
[Unit]
Description=Hermes — load secrets from Secret Manager into /run/hermes/env
After=network-online.target
Wants=network-online.target
Before=docker-compose-hermes.service
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/local/bin/load-secrets.sh
RemainAfterExit=yes
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target
